# usb-modem-huawei-e175x/files/etc/uci-defaults/usb-modem-huawei-e175x
#!/bin/sh

# Copyright (C) 2010 ePoint Systems Ltd.
# Author: Daniel A. Nagy

uci batch <<-EOF
    set network.mobile=interface
    set network.mobile.device=/dev/ttyUSB0
    commit network
EOF
