# flexibity-sensors/files/tunslip6.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 OpenWrt.org
START=50

SSD=start-stop-daemon
NAME=tunslip6
PIDF=/var/run/$NAME.pid
PROG=/usr/bin/$NAME
ADDR=aaaa::1/64

start() {
	mc1322x-reset.sh && sleep 3 && $SSD -S -m -p $PIDF -q -x $PROG -- -B 115200 -s /dev/ttyS1 -t wpan0 $ADDR &
}

stop() {
	$SSD -K -p $PIDF
}

